# crontab
0 10 * * * python manage.py get_playlist_updates --send-notifications
